# .config/systemd/user/dns-update.service
[Unit]
Description=Dynamic DNS update

[Service]
ExecStart=/home/cagprado/bin/dns-update.sh
