# .config/systemd/user/dns-update.service
[Unit]
Description=Dynamic DNS update

[Service]
ExecStart=/home/cagprado/.config/systemd/user/scripts/dns-update.zsh
