# .config/systemd/user/internet-check.service
[Unit]
Description=Internet connection check

[Service]
ExecStart=/home/cagprado/bin/internet-check.sh

[Install]
WantedBy=basic.target
