# .config/systemd/user/mail.timer
[Unit]
Description=Mail Update Timer

[Timer]
OnBootSec=0
OnUnitActiveSec=4min

[Install]
WantedBy=timers.target
