# .config/systemd/user/mail.service
[Unit]
Description=Notmuch mail update

[Service]
ExecStart=/home/cagprado/bin/mail-update.sh
