# .config/systemd/user/mail.service
[Unit]
Description=Notmuch mail update

[Service]
ExecStart=/bin/zsh -c '~/bin/zsh/mail-update'
