# .config/systemd/user/network-services-mredson.target
[Unit]
Description=Network services
#Wants=dns-update.timer
