# .config/systemd/user/network-services-other.target
[Unit]
Description=Network services
Wants=dropbox.service
Wants=mail.timer
