# .config/systemd/user/network-services-common.target
[Unit]
Description=Network services (Common)
#Wants=mail.timer
