# .config/systemd/user/network-services.target
[Unit]
Description=Network services
Wants=mail.timer
Wants=contacts.timer
