# .config/systemd/user/vdirsyncer.service
[Unit]
Description=vdirsyncer update

[Service]
ExecStart=/bin/zsh -c '~/bin/zsh/contacts-update'
