# etc/system/etc/profile.d/texlive.sh
export PATH=/usr/local/texlive/bin/x86_64-linux:${PATH}
#export MANPATH=/usr/local/texlive/texmf-dist/doc/man:${MANPATH}
#export INFOPATH=/usr/local/texlive/texmf-dist/doc/info:${INFOPATH}
